# alicia-brain/scripts/com.hygge.white-rabbit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hygge.white-rabbit</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/eduardobonilla/.volta/bin/node</string>
    <string>/Users/eduardobonilla/Desktop/ALICE/alicia-brain/scripts/scrape.js</string>
    <string>all</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/eduardobonilla/Desktop/ALICE/alicia-brain</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>BRAIN_URL</key>
    <string>https://aliceai.bam.pe</string>
    <key>MARKET_REFRESH_TOKEN</key>
    <string>white-rabbit</string>
    <key>NODE_PATH</key>
    <string>/usr/local/lib/node_modules</string>
  </dict>

  <!-- Run every 6 hours -->
  <key>StartInterval</key>
  <integer>21600</integer>

  <key>StandardOutPath</key>
  <string>/Users/eduardobonilla/Library/Logs/white-rabbit.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/eduardobonilla/Library/Logs/white-rabbit.error.log</string>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
